# deploy/com.dealerplatform.pipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- Runs the Credit Pipeline lead poller every 5 minutes as a user LaunchAgent.
     Copy to ~/Library/LaunchAgents/ and `launchctl load` it.
     Each run: reject stale unsent leads, enrich contact from tbl_ownership into
     match_phone/match_email, then send (newest-first). The 5-minute cadence gives
     the ownership match time to complete before the send it gates.
     launchd will not start a second copy while one is still running, so runs
     never overlap. The poller checks the global flow switch (platform admin) and
     exits immediately when it is OFF, so it is cheap to leave loaded. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dealerplatform.pipeline</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/markrankin/claude/credit-estimator/run_poller.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>300</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/markrankin/claude/deploy/pipeline.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/markrankin/claude/deploy/pipeline.err.log</string>
</dict>
</plist>
