# deploy/com.dealerplatform.scraper.adesa.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- Runs the Adesa scraper every 48 hours as a user LaunchAgent. Load only after
     the one-time login:
       cp com.dealerplatform.scraper.adesa.plist ~/Library/LaunchAgents/
       launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.dealerplatform.scraper.adesa.plist
     Runs headful in the logged-in session; launchd won't start a second copy while
     one is still running. Adjust StartInterval to change the cadence. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dealerplatform.scraper.adesa</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/markrankin/claude/scrapers/run.sh</string>
        <string>adesa</string>
    </array>
    <key>StartInterval</key>
    <integer>172800</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/markrankin/claude/deploy/scraper.adesa.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/markrankin/claude/deploy/scraper.adesa.err.log</string>
</dict>
</plist>
