# deploy/com.dealerplatform.pipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- Runs the Credit Pipeline lead poller as a persistent user LaunchAgent.
     run_poller.sh is a supervised loop: it runs the poller once per minute in a
     fresh python process, forever. KeepAlive=true means launchd relaunches the
     script if it ever exits (crash, kill, logout/login, reboot) — so the send
     schedule no longer depends on launchd's StartInterval timer, which stalled
     unreliably on this host and silently killed the send about once an hour.
     Each tick is time-capped inside run_poller.sh so a stuck run can't freeze the
     loop. The poller checks the global flow switch (platform admin) and exits
     immediately when it is OFF, so it is cheap to leave running. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dealerplatform.pipeline</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/markrankin/claude/credit-estimator/run_poller.sh</string>
    </array>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>/Users/markrankin/claude/deploy/pipeline.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/markrankin/claude/deploy/pipeline.err.log</string>
</dict>
</plist>
